# site-cookbooks/gitlab/templates/default/rc.local.erb
#!/bin/sh -e
#
# rc.local
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.
#
# In order to enable or disable this script just change the execution
# bits.
#
# By default this script does nothing.


test -h /vagrant && rm -fv /vagrant
test -d /vagrant && rmdir -v /vagrant
test -d /opt/ruby -a ! -d /opt/vagrant_ruby && ln -sf /opt/ruby /opt/vagrant_ruby

exit 0
